# autobackup/backup_logic.timer
[Unit]
Description=Run backup_logic.service weekly

[Timer]
OnActiveSec=1m
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
