# autobackup/binlog_backup.timer
[Unit]
Description=Run binlog_backup.service every hour

[Timer]
OnActiveSec=1m
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
